# rc.local
#!/bin/sh -e
/usr/sbin/coursework
exit 0
